# dockerfiles/php/supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/www/storage/logs/supervisord_%(ENV_PROCESS)s.log
pidfile=/var/www/storage/logs/supervisord_%(ENV_PROCESS)s.pid

[program:php]
directory=/var/www/
command=/bin/sh -c "exec %(ENV_COMMAND)s"
autostart=true
autorestart=false
startretries: 3 
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:horizon]
command=php /var/www/artisan horizon
process_name=horizon-%(process_num)02d
directory=/var/www
autostart=true
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PROCESS="horizon"